# _/cron/com.user.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://apple.com">
<plist version="1.0">
<!-- https://gist.github.com/yoginsurya/f7ee94550d6d421ea2675fa4d6aee281 -->
<dict>
    <key>Label</key>
    <string>com.user.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Ghostty.app/Contents/MacOS/ghostty</string>
        <string>-e</string>
        <string>/Users/jms/dev/_/cron/daily.sh</string>
        <string>bug-archive</string>
    </array>

    <!-- Help the script detect whether launchd is running it. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>CRON_DIR</key>
        <string>/Users/jms/_/cron</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WakeSystem</key>
    <true/>

    <!-- Essential so Ghostty can connect to macOS GUI session -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>StandardOutPath</key>
    <string>/Users/jms/_/cron/logs/daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jms/_/cron/logs/daily.err</string>
</dict>
</plist>

